# debian/cron.d
# Every minute, do sympa actions
* * * * *	   root		/usr/lib/alternc/update_sympa.sh
