# debian/cron.d
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin
# Every minute, do sympa actions
* * * * *	   root		/usr/lib/alternc/update_sympa.sh
